# deploy/orderbook.service
[Unit]
Description=Order book matching engine and web front end
After=network.target

[Service]
Environment=JAVA_HOME=/opt/jdk25
Environment=JAVA_OPTS=-Xmx256m
Environment=PORT=8080
# Box-local egress config (KAFKA_BOOTSTRAP_SERVERS etc.); optional - absent means egress off.
EnvironmentFile=-/etc/orderbook/egress.env
ExecStart=/home/ubuntu/orderbook/bin/orderbook
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
User=ubuntu

[Install]
WantedBy=multi-user.target
